# data/com.pixar.xoloserver.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>KeepAlive</key>
  <true/>
  <key>Label</key>
  <string>com.pixar.xoloserver</string>
  <key>EnvironmentVariables</key>
  <dict>
    <key>GEM_PATH</key>
    <string>/usr/local/pixar/ruby/aliases/default/lib/ruby/gems:/usr/local/pixar/lib/ruby/gems/site_ruby:/usr/local/pixar/lib/ruby/gems/pixar:/usr/local/pixar/lib/ruby/gems/2.7.0</string>
  </dict>
  <key>ProgramArguments</key>
  <array>
    <string>/usr/local/bin/xoloserver</string>
  </array>
  <key>RunAtLoad</key>
  <true/>
  <key>StandardOutPath</key>
  <string>/Library/Application Support/xoloserver/logs/xoloserver.log</string>
  <key>StandardErrorPath</key>
  <string>/Library/Application Support/xoloserver/logs/xoloserver.log</string>
</dict>
</plist>
